# infra/systemd/demin-auto-approve.service
[Unit]
Description=DEMIN auto_approve -- aprueba drafts auto cuando hitl_mode=false (Sprint 6 Fase 3)
Documentation=https://github.com/demin-group/demin-system
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=demin
Group=demin
WorkingDirectory=/home/demin/demin-system/apps/workers
Environment=ENV=prod
Environment=PATH=/home/demin/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/home/demin/.local/bin/uv run python -m outreach.auto_approve --env prod

TimeoutStartSec=300
Nice=10
IOSchedulingClass=idle

StandardOutput=journal
StandardError=journal
SyslogIdentifier=demin-auto-approve

Restart=no

[Install]
WantedBy=multi-user.target
