# infra/systemd/demin-handle-actions.service
[Unit]
Description=DEMIN handle_actions -- ejecuta accion por categoria de reply (Fase 3)
Documentation=https://github.com/demin-group/demin-system
After=network-online.target demin-classify-replies.service
Wants=network-online.target

[Service]
Type=oneshot
User=demin
Group=demin
WorkingDirectory=/home/demin/demin-system/apps/workers
Environment=ENV=prod
Environment=PATH=/home/demin/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/home/demin/.local/bin/uv run python -m replies.handle_actions --env prod

TimeoutStartSec=300
Nice=10
IOSchedulingClass=idle

StandardOutput=journal
StandardError=journal
SyslogIdentifier=demin-handle-actions

Restart=no

[Install]
WantedBy=multi-user.target
