# infra/systemd/demin-followups.service
[Unit]
Description=DEMIN follow_ups — programa D+4 (reframe) y D+10 (closing) en cola HITL
Documentation=https://github.com/demin-group/demin-system
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=demin
Group=demin
WorkingDirectory=/home/demin/demin-system/apps/workers
Environment=ENV=prod
Environment=PATH=/home/demin/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/home/demin/.local/bin/uv run python -m outreach.follow_ups --env prod

TimeoutStartSec=600
Nice=10
IOSchedulingClass=idle

StandardOutput=journal
StandardError=journal
SyslogIdentifier=demin-followups

Restart=no

[Install]
WantedBy=multi-user.target
