# infra/systemd/demin-auto-pause.service
[Unit]
Description=DEMIN auto_pause — vigila bounce >2% / spam >0.1% rolling 7d y pausa mailbox si excede
Documentation=https://github.com/demin-group/demin-system
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=demin
Group=demin
WorkingDirectory=/home/demin/demin-system/apps/workers
Environment=ENV=prod
Environment=PATH=/home/demin/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/home/demin/.local/bin/uv run python -m outreach.auto_pause --env prod

TimeoutStartSec=300
Nice=10
IOSchedulingClass=idle

StandardOutput=journal
StandardError=journal
SyslogIdentifier=demin-auto-pause

Restart=no

[Install]
WantedBy=multi-user.target
